# version-checker/crontab
0 */8 * * * /bin/zsh /Users/nrjdalal/.scripts/version-checker.sh > /Users/nrjdalal/.scripts/version-checker.log 2>&1
